# scripts/service-templates/expertise-api.service.tmpl
[Unit]
Description=Agent Expertise API
Documentation=https://github.com/TheSemicolon/agent-expertise-api
# Ordering only — user units cannot Require= system targets.
After=network-online.target

[Service]
Type=notify
NotifyAccess=main
WatchdogSec=30
WorkingDirectory=@@WORKDIR@@
ExecStart=@@WRAPPER@@

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Environment file (optional — install.sh creates a chmod-600 stub)
EnvironmentFile=-@@SECRETS_FILE@@

# Hardening (user-unit subset; some are no-ops at user scope but documented
# for system-unit promotion).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native

# DO NOT enable MemoryDenyWriteExecute — RyuJIT writes executable pages.
# Enabling crashes the .NET process at first JIT.
MemoryDenyWriteExecute=false

# Allow log dir writes (user-unit honors this; explicit for clarity)
ReadWritePaths=@@LOG_DIR@@

[Install]
# User unit — default.target, NOT multi-user.target.
WantedBy=default.target
